# etc/init.d/ttpd-orig.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          @PACKAGE_NAME@
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Init script for the TTPD server
# Description:       This init.d script is used to start @PACKAGE_NAME@.
#                    It basically just calls @PACKAGE_NAME@ctl.
### END INIT INFO

# Author:	Martin Thorsen Ranang <mtr@ranang.org>.

PATH=$HOME/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=@PACKAGE_NAME@
DAEMON_CONTROLLER=@PACKAGE_NAME@ctl
NAME=@PACKAGE_NAME@
DESC="the TUC Transfer Protocol Daemon (TTPD)"
LOG=/var/log/@PACKAGE_NAME@/@PACKAGE_NAME@.log

# Gracefully exit if the package has been removed.
#test -x $DAEMON || exit 0

# Read config file if it is present.
#[ -r /etc/default/$NAME ] && . /etc/default/$NAME

set -e

FLAGS="-Ldebug --tuc-threads=1"

get_@PACKAGE_NAME@_pid()
{
    PID=`grep "@PACKAGE_NAME@ INFO PID =" $LOG |tail -1 |cut -f 2 -d '=' |cut -f 1 -d ','`;
    return
}

case "$1" in
  start)
	echo -n "Starting $DESC: $NAME"
	$DAEMON #$FLAGS #start
	echo "."
	;;
  no-daemon)
	echo -n "Starting $DESC: $NAME"
	$DAEMON #$FLAGS --no-daemon &
	echo "."
	;;
  stop)
	echo -n "Stopping $DESC: $NAME"
	get_@PACKAGE_NAME@_pid
	kill $PID
	echo "."
	;;
  force-stop)
	echo -n "Force-stopping $DESC: $NAME"
	get_@PACKAGE_NAME@_pid
	kill -9 $PID
	echo "."
	;;
  restart)
	echo -n "Restarting $DESC configuration..."
	get_@PACKAGE_NAME@_pid
	kill -HUP $PID
	echo "done."
	;;
  store-state)
	echo -n "Storing scheduler state in $DESC: $NAME"
	get_@PACKAGE_NAME@_pid
	kill -USR1 $PID
	echo "."
	;;
  *)
	N=/etc/init.d/$NAME
	echo "Usage: $N {start|stop|restart|store-state}" >&2
	exit 1
	;;
esac

exit 0
